# src/3rdparty/applications/qx/services/QX.service
[Translation]
File=QtopiaServices
Context=QX
[Service]
Actions = "appSwitch()"
Icon = qx/qx
Name[]=QX
[appSwitch()]
Icon=qx/qx
Name[]="Switch X application"
